# roles/flask/files/flask_app.conf
[program:basic_app]
command= /usr/bin/python /home/max/flask_app/hello.py
autostart=true
autorestart=true
